# deploy/hamsvic-celery.service
# ==============================================================================
# Systemd Service for HAMSVIC Celery Worker
# ==============================================================================
# Place at: /etc/systemd/system/hamsvic-celery.service

[Unit]
Description=HAMSVIC Celery Worker
After=network.target redis.service
Requires=network.target

[Service]
Type=forking
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/hamsvic
Environment="PATH=/home/ubuntu/hamsvic/venv/bin"
EnvironmentFile=/home/ubuntu/hamsvic/.env
ExecStart=/home/ubuntu/hamsvic/venv/bin/celery \
    -A estimate_site worker \
    --loglevel=info \
    --concurrency=4 \
    --pidfile=/run/celery/hamsvic-worker.pid \
    --logfile=/var/log/celery/hamsvic-worker.log
ExecStop=/bin/kill -s TERM $MAINPID
ExecReload=/bin/kill -s HUP $MAINPID
PIDFile=/run/celery/hamsvic-worker.pid
Restart=always
RestartSec=10

# Create runtime directory
RuntimeDirectory=celery

[Install]
WantedBy=multi-user.target
